# service/appman_api.service
[Unit]
Description=The API for appman
[Service]
ExecStart=appman_api.sh
WorkingDirectory=/home/helauren/.local/appman/api
Restart=always
User=helauren
Group=helauren
[Install]
WantedBy=default.target
